# .dockers/supervisord/supervisord-mailcatcher.conf
[program:mailcatcher]
command=mailcatcher -f --http-ip=0.0.0.0
stdout_events_enabled=true
stderr_events_enabled=true
priority=300
