# .dockers/supervisord/supervisord-mysqld.conf
[program:mysqld]
command=/usr/bin/mysqld_safe
numprocs=1
autostart=true
autorestart=false
priority=100