# .dockers/supervisord/supervisord-apache2.conf
[program:apache2]
command=apache2ctl -DFOREGROUND
numprocs=1
autostart=true
autorestart=true