# .dockers/supervisord/supervisord-cron.conf
[program:syslog-ng]
command=/usr/sbin/syslog-ng --foreground  --no-caps
redirect_stderr=true
autostart=true
stdout_syslog=true