# .dockers/supervisord/supervisord-syslog.conf
[program:cron]
command=/usr/sbin/cron -f
numprocs=1
autostart=true
autorestart=false
priority=200